# com.bing.wallpaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- 唯一标识 -->
    <key>Label</key>
    <string>com.bing.wallpaper</string>

    <!-- 程序和参数 -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bing-wallpaper-macos</string>
        <string>--auto</string>
    </array>

    <!-- 开机自动运行 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 每天早上 8 点运行一次 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 日志输出 -->
    <key>StandardOutPath</key>
    <string>/tmp/bing-wallpaper.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bing-wallpaper.err</string>
</dict>
</plist>

